# profile.d/keyboard.sh
# Caps lock -> control
setxkbmap -option ctrl:nocaps
# Natural scrolling
synclient VertScrollDelta=-111
synclient HorizScrollDelta=-111
